# artifact_evaluation/siloz.service
[Unit]
Description=Siloz hpage systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /root/reserve_hugepages.sh

[Install]
WantedBy=multi-user.target

